# spatial-information-subsystem/com-collection-app/weather-information-collector/config/systemd/weather-information-collector.service
[Unit]
Description=weather-information-collector
After=network.target

[Service]
User=digiline
ExecStart=/opt/weather-information-collector/bin/bootstrap.sh
SuccessExitStatus=143
Restart=always

[Install]
WantedBy=multi-user.target