# spatial-information-subsystem/com-collection-app/target-information-collector/config/systemd/target-information-collector.service
[Unit]
Description=target-information-collector
After=network.target

[Service]
User=digiline
ExecStart=/opt/target-information-collector/bin/bootstrap.sh
SuccessExitStatus=143
Restart=always

[Install]
WantedBy=multi-user.target